# deploy/systemd/tenmatch-watchdog.service
[Unit]
Description=TenMatch local liveness and memory-pressure watchdog
Documentation=file:/srv/tenmatch/app/deploy/README.md
After=tenmatch.service
ConditionPathExists=/sys/fs/cgroup/cgroup.controllers

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/libexec/tenmatch-watchdog
TimeoutStartSec=120s

RuntimeDirectory=tenmatch-watchdog
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tenmatch-watchdog
UMask=0077
Nice=10
IOSchedulingClass=idle

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET
